# files/etc/uci-defaults/99-custom
#!/bin/sh
# Oasisic OpenWrt — 首次启动补丁
# 创建 .oasisic-firstboot 标记，让 index.html 跳转到设置向导

# ---- 创建首次启动标记（设置向导用）----
touch /etc/.oasisic-firstboot

# ---- 设置 LuCI 语言 & 主题 & 资源路径（bootstrap 默认）----
uci set luci.languages='internal'
uci set luci.languages.zh_cn='简体中文 (Simplified Chinese)'
uci set luci.main.lang='zh_cn'
uci set luci.main.mediaurlbase='/luci-static/bootstrap'
uci set luci.main.resourcebase='/luci-static/resources'

# ---- 网络诊断默认地址（DNSPod）----
uci set luci.diag.dns='119.29.29.29'
uci set luci.diag.ping='119.29.29.29'
uci set luci.diag.route='119.29.29.29'

uci commit luci

# ---- 默认 DHCP 客户端（无静态 IP）----
uci set network.lan.proto='dhcp'

# ---- IPv6 DNS（国内常用）----
uci add_list network.lan.dns='2400:3200::1'
uci add_list network.lan.dns='2402:4e00::6666'

# ---- IPv6 委派默认关闭（旁路网关不向下分发 IPv6）----
uci set network.lan.delegate='0'

# ---- IPv6 中继默认关闭 ----
uci set dhcp.lan.dynamicdhcp='0'
uci set dhcp.lan.dhcpv6='disabled'
uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.ndp='disabled'

uci commit

# ---- 确保 QEMU Guest Agent 启动 ----
if [ -x /etc/init.d/qemu-ga ]; then
    /etc/init.d/qemu-ga enable
    /etc/init.d/qemu-ga start
fi

# ---- 清理自身 ----
rm /etc/uci-defaults/99-custom
exit 0